# acestream.conf
[program:acestreamengine]
command=/usr/bin/acestreamengine --bind-all --client-console
